# templates/traefik.service.j2
[Unit]
After=network.service

[Service]
WorkingDirectory=
ExecStart={{ traefik_dir }}/traefik/traefik
Restart=always

[Install]
WantedBy=default.target
