# templates/traefik.service.j2
[Unit]
After=network.service

[Service]
WorkingDirectory=
ExecStart=traefik --configFile=/etc/traefik/static.yaml
Restart=always

[Install]
WantedBy=default.target
